# node_api.ini
[program:node_api]
command=/var/communication-backend/src/node_api.py 8080
autostart=true
autorestart=true
startretries=3
