# node_api.ini
[program:node_api]
priority=1
directory=/var/communication-backend
command=/var/communication-backend/venv/bin/python src/node_api.py -b %(ENV_BOOTSTRAP_NODE)s %(ENV_BOOTSTRAP_PORT)s
autostart=true
autorestart=true
startretries=3
