# node_api.ini
[program:node_api]
directory=/var/communication-backend
command=/var/communication-backend/venv/bin/python src/node_api.py 5000
autostart=true
autorestart=true
startretries=3
