# postgres.ini
[program:postgresql]
priority=100
directory=/var/lib/postgresql/data
command=/usr/bin/postgres
user=postgres
autorestart=true
autostart=true
startretries=3
