# register_service.ini
[program:register_service]
command=curl -i -H "Content-Type: application/json" -X POST -d '{"service_name":"%(ENV_SERVICE_NAME)s"}' http://localhost:9001/register
startsecs=0
autostart=false
autorestart=false
