# register_service.ini
[program:register_service]
priority=99999
command=sh -c 'sleep 3 && /var/communication-backend/register_service.sh'
startsecs=0
autostart=true
autorestart=false
stderr_logfile = /var/log/register-stderr.log
stdout_logfile = /var/log/register-stdout.log
