# register_service.ini
[program:register_service]
command=/var/communication-backend/register_service.sh
startsecs=0
autostart=true
autorestart=false
stderr_logfile = /var/log/register-stderr.log
stdout_logfile = /var/log/register-stdout.log
